# templates/etc/update-motd.d/motd.txt.j2
#!/bin/sh
#
# Managed by Ansible
#
printf "\n%s" '{{ motd_content }}'
